# Actividad 4/saludar.service
[Unit]
Description=Saludar y mostrar fecha
[Service]
Type=simple
ExecStart= saludo.sh
[Install]
WantedBy=multi-user.target